# rpi_image_scripts/bookworm/service/de_camera_rpi_cam.service
[Unit]
Description=Drone Engage Camera with Initialization of rpicamera-vid pipeline Ver 2.0
After=network.target multi-user.target
Documentation=https://droneengage.com


[Service]
Type=simple
WorkingDirectory=/home/pi/
ExecStart=/home/pi/scripts/wrapper/camera_manager_wrapper -c 
Restart=on-failure
RestartSec=5s # Wait 5 seconds before attempting a restart

StandardOutput=journal
StandardError=journal
SyslogIdentifier=camera_wrapper

# You might need to adjust user/group if camera access is restricted
# User=pi
# Group=video

[Install]
WantedBy=multi-user.target

